# docker/build/php-cli-7.1/supervisor/conf.d/cabinet.socket.conf
[program:cabinet-socket]
command=nodejs /var/www/cabinet/app/nodejs/server.js
user=user
stdout_logfile=/var/log/supervisor/cabinet.log
stderr_logfile=/var/log/supervisor/cabinet.error.log