# docker/build/php-cli-7.1/supervisor/conf.d/rabbit-mq.conf
[program:rabbit-avto-servisy-group-work-all]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/group-work/all --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-group-work]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/group-work/group-work --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-group]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/group-work/group --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-work]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/group-work/work --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-discount]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/discount/run --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-delete]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/discount/delete --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-discount-all]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/discount/all --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-form-enable]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/form/enable --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-form-disable]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/form/disable --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-form-delete]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/form/delete
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-cache]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/cache/run --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-comment]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/comment/run --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-comment-delete]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/comment/delete --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-comment-all]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/comment/all --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-photo]
command=php /home/srv/vse-avtoservisy.ru/app/yii amqp/photo/run --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log