# docker/build/php-cli-7.1/supervisor/conf.d/rabbit-mq.conf
[program:rabbit-avto-servisy-group-work-all]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/group-work/all --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-group-work]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/group-work/group-work --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-group]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/group-work/group --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-work]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/group-work/work --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-discount]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/discount/run --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-form]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/form/run --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

[program:rabbit-avto-servisy-cache]
command=php /var/www/vse-avtoservisy.ru/app/yii amqp/cache/run --log
user=user
stdout_logfile=/var/log/supervisor/rabbit-avto-servisy.log
stderr_logfile=/var/log/supervisor/rabbit-avto-servisy.error.log

