# docker/build/php-cli-7.3/supervisor/conf.d/ttr-back.conf
[program:laravel-worker-default]
command=php /var/www/ttr-back/artisan queue:work --queue=default --sleep=10 --timeout=30
user=user
stdout_logfile=/var/log/supervisor/ttr-back.log
stderr_logfile=/var/log/supervisor/ttr-back.error.log

[program:laravel-worker-emails]
command=php /var/www/ttr-back/artisan queue:work --queue=emails --sleep=10 --timeout=31 --tries=3
user=user
stdout_logfile=/var/log/supervisor/ttr-back.log
stderr_logfile=/var/log/supervisor/ttr-back.error.log

[program:laravel-worker-presentations]
command=php /var/www/ttr-back/artisan queue:work --queue=presentations --sleep=1 --timeout=360
user=user
stdout_logfile=/var/log/supervisor/ttr-back.log
stderr_logfile=/var/log/supervisor/ttr-back.error.log